# daemon/tk-daemon.socket
[Unit]
Description=Timekeeping socket

[Socket]
ListenStream = %t/tk.sock

[Install]
WantedBy = sockets.target
